# home/crontabs/usercrontab
@daily /$HOME/sys-info.sh >> /home/pi/afg-ericsson.log
@reboot . /$HOME/.bash_ericsson; /home/pi/command-receiver.sh >> /home/pi/afg-ericsson.log
