# home/crontabs/usercrontab
@daily . /etc/profile.d/afg_ericsson.sh; . /etc/profile.d/raspberry_id.sh; "$HOME"/sys-info.sh
@reboot . /etc/profile.d/afg_ericsson.sh; . /etc/profile.d/raspberry_id.sh; "$HOME"/command-receiver.sh > "$HOME"/cmd-log.txt
0 16-19 * * * sudo "$HOME"/main-updater.sh
