# home/crontabs/usercrontab
@daily . "$HOME"/.bash_afg_ericsson; "$HOME"/sys-info.sh
@reboot . "$HOME"/.bash_afg_ericsson; "$HOME"/command-receiver.sh
0 16-19 * * * sudo "$HOME"/main-updater.sh
