# home/crontabs/usercrontab
@daily "$HOME"/sys-info.sh
@reboot "$HOME"/command-receiver.sh
0 16-19 * * * sudo "$HOME"/main-updater.sh
